# .docker/dev/supervisor/supervisor.ini
[supervisord]
nodaemon=true

[program:adicionar-dados-na-base]
command=cd /docker/ && teste.sh
startsecs = 0
autorestart = false
startretries = 1
